# DataCompute/python/crontab/transfer_history_log_to_cdh4.sh
#!/bin/bash

for id in 14 15 16 17 18
do
    ./transfer_minik_log_to_cdh4.sh 2016-09-$id
done
